# config/hostapd.service
[Unit]
Description=Wi-Fi AP (hostapd)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/sbin/hostapd -B /etc/hostapd/hostapd.conf
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
